# script/bme.service
[Unit]
Description = BME control

[Service]
ExecStart =