# script/bme.service
[Unit]
Description=BME control

[Service]
ExecStart=/home/pi/git/RunBme/script/restart_bme.sh
Restart=always

[Install]
WantedBy=multi-user.target